# etc/supervisord/conf.d/flask_app_sample.conf
[program:flask_app_sample]
command = gunicorn flask_app_sample:app --bind 0.0.0.0:2105 --workers 4
user = nobody
autostart = true
autorestart = true
